# web-wallpaper-api.service
[Unit]
Description=Web Wallpaper Flask API Service
After=network.target

[Service]
Type=simple
User=%INSTALL_USER%
WorkingDirectory=%INSTALL_DIR%
Environment="PATH=%INSTALL_DIR%/.venv/bin:/usr/local/bin:/usr/bin:/bin"
# Bind to localhost only for security - no external access
ExecStart=%INSTALL_DIR%/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
